# GTM/etc/xinetd.d/vista-hl7
service vista-hl7
{
  port = 5001
  socket_type = stream
  protocol = tcp
  type = UNLISTED
  user = foia
  server = /home/foia/bin/hl7.sh
  wait = no
  disable = no
  per_source = UNLIMITED
  instances = UNLIMITED
}
